# deb_build/DEBIAN/prerm

#!/usr/bin/env bash
set -euo pipefail

# Pre-removal script for the nideflanders package
# Stops and disables the systemd unit if present.

if command -v systemctl >/dev/null 2>&1; then
  if systemctl is-active --quiet nideflanders.service; then
    systemctl stop nideflanders.service || true
  fi
  if systemctl is-enabled --quiet nideflanders.service; then
    systemctl disable nideflanders.service || true
  fi
fi

exit 0
#!/usr/bin/env bash
set -euo pipefail

# Pre-removal script for nideflanders .deb
if systemctl is-active --quiet nideflanders.service; then
  systemctl stop nideflanders.service || true
fi
if systemctl is-enabled --quiet nideflanders.service; then
  systemctl disable nideflanders.service || true
fi

exit 0
